# services/robotSystem.target
[Unit]
Description=Robot system
Requires=robot@hardware.service robot@camera.service robot@ros_server.service robot@user_server.service robot@vision.service
After=network.target