# services/robotSystem.target
[Unit]
Description=Robot system
Requires=robot@hardware.service
After=network.target